# solis-ea-modbus/rootfs/etc/services.d/solis-ea/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Solis EA Modbus Control — S6 service finish script
# ==============================================================================

declare exit_code=${1}

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.warning "Controller exited with code ${exit_code} — S6 will restart"
fi
