# solis-ea-modbus/rootfs/etc/services.d/solis-ea/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Solis EA Modbus Control — S6 service run script
# ==============================================================================

# Read add-on options
GATEWAY_HOST=$(bashio::config 'gateway_host')
GATEWAY_PORT=$(bashio::config 'gateway_port')
PROTOCOL=$(bashio::config 'protocol')
SLAVE_ADDRESS=$(bashio::config 'slave_address')
POLL_INTERVAL=$(bashio::config 'poll_interval')
INTER_FRAME_MS=$(bashio::config 'inter_frame_ms')
ENABLE_CONTROL=$(bashio::config 'enable_control')
MAX_CURRENT_A=$(bashio::config 'max_current_a')
EXPORT_LIMIT_W=$(bashio::config 'export_limit_w')
DEBUG=$(bashio::config 'debug')

# MQTT: prefer add-on options, fall back to HA MQTT service discovery
MQTT_HOST=$(bashio::config 'mqtt_host')
MQTT_PORT=$(bashio::config 'mqtt_port')
MQTT_USER=$(bashio::config 'mqtt_user')
MQTT_PASS=$(bashio::config 'mqtt_pass')

if bashio::var.is_empty "${MQTT_HOST}"; then
    if bashio::services.available "mqtt"; then
        MQTT_HOST=$(bashio::services mqtt "host")
        MQTT_PORT=$(bashio::services mqtt "port")
        MQTT_USER=$(bashio::services mqtt "username")
        MQTT_PASS=$(bashio::services mqtt "password")
        bashio::log.info "Using HA MQTT service: ${MQTT_HOST}:${MQTT_PORT}"
    else
        bashio::log.warning "No MQTT service available — publishing disabled"
    fi
fi

# Export for Python script
export GATEWAY_HOST
export GATEWAY_PORT
export PROTOCOL
export SLAVE_ADDRESS
export MQTT_HOST
export MQTT_PORT
export MQTT_USER
export MQTT_PASS
export POLL_INTERVAL
export INTER_FRAME_MS
export ENABLE_CONTROL
export MAX_CURRENT_A
export EXPORT_LIMIT_W
export LOG_DIR="/config"

if bashio::var.true "${DEBUG}"; then
    export DEBUG="true"
fi

bashio::log.info "Starting Solis EA Modbus Control"
bashio::log.info "  Gateway: ${GATEWAY_HOST}:${GATEWAY_PORT} (${PROTOCOL}), unit ${SLAVE_ADDRESS}"
bashio::log.info "  MQTT: ${MQTT_HOST}:${MQTT_PORT}"
bashio::log.info "  Control: ${ENABLE_CONTROL} (max current ${MAX_CURRENT_A} A, export <= ${EXPORT_LIMIT_W} W)"
bashio::log.info "  Poll: ${POLL_INTERVAL}s, inter-frame ${INTER_FRAME_MS}ms"

mkdir -p "${LOG_DIR}"

# Run the controller — S6 will restart if it exits
exec python3 /opt/solis/solis_ea.py
